# systemd.confs/jack-healthcheck.timer
[Unit]
Description=Run JACK monitor every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=30s

[Install]
WantedBy=timers.target
